# docker/supervisord.conf
[supervisord]
nodaemon=true
loglevel=info
logfile=/dev/stdout
logfile_maxbytes=0
pidfile=/var/run/supervisord.pid
user=root

# Virtual framebuffer — the X server every other component talks to.
[program:xvfb]
command=/usr/bin/Xvfb :99 -screen 0 1280x800x24 -nolisten tcp -ac
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# Bridges the Xvfb display to a VNC server on 5900 (loopback only — noVNC
# proxies it out). `-shared` lets multiple noVNC tabs view at once, which is
# handy if you want to debug from two screens.
[program:x11vnc]
command=/usr/bin/x11vnc -display :99 -nopw -listen 127.0.0.1 -rfbport 5900 -shared -forever -quiet
autostart=true
autorestart=true
priority=20
startsecs=2
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# noVNC is the static frontend; websockify is the actual TCP-to-WebSocket
# bridge. Together they turn http://host:6080 into a working VNC client.
[program:novnc]
command=/usr/bin/websockify --web=/usr/share/novnc 6080 127.0.0.1:5900
autostart=true
autorestart=true
priority=30
startsecs=2
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# Agent web UI on :8000 — the noVNC iframe inside it points to :6080. The
# server starts even without an API key; Settings.load() only fires when
# the user hits Run, so a missing key produces an in-UI error rather than
# a crashed container.
[program:web]
command=/usr/bin/python -m agent serve --host 0.0.0.0 --port 8000
directory=/app
autostart=true
autorestart=true
priority=40
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
